# mountnas.conf
# mount nas forlder when WIFI is ready
# sshfsm must be in /bin and this script in /etc/init/

start on net-device-up wlan0

exec  sshfsm mount
